# scripts/systemd/myanmar-api-v2.service
[Unit]
Description=TOPIK Myanmar FastAPI (apps/api)
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/myanmar-v2/apps/api
EnvironmentFile=/opt/myanmar-v2/apps/api/.env
ExecStart=/opt/myanmar-v2/apps/api/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000 --workers 2
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
